# vagrant/supervisor/worker-400667.conf
[program:worker-progress-01]
command=php8.2 /home/forge/progress/artisan queue:work beanstalkd --sleep=0.5 --daemon --quiet --timeout=30 --memory=128 --tries=1 --env=local --queue="default"

process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=vagrabt
numprocs=2
stdout_logfile=/home/vagrant/.forge/worker-progress-01.log
